# zelotoj.route/upload-files/etc/crontabs/root
0 5 * * * /root/gfwlist2dnsmasq.py && mv /root/dnsmasq_list.conf /etc/dnsmasq.d/dnsmasq_list.conf
