# zelotoj.route/upload-files/etc/crontabs/root
0 5 * * * /usr/share/gfwlist/gfwlist_update.sh
